# .config/autostart/Picom.desktop
[Desktop Entry]
Exec=picom
Name=Picom
Type=Application
Version=1.0
